# mkosi.extra/etc/systemd/system/rsyslog-hostname-reload.path
[Unit]
Description=Watch for hostname changes and reload rsyslog

[Path]
PathChanged=/etc/hostname
PathChanged=/run/systemd/netif/leases

[Install]
WantedBy=multi-user.target
