# mkosi.extra/etc/systemd/system/rsyslog-hostname-reload.service
[Unit]
Description=Restart rsyslog to pick up new hostname

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart rsyslog.service
