# mkosi.extra/etc/systemd/system/tink-agent.service
[Unit]
Description=Tinkerbell Agent
After=tink-agent-setup.service containerd.service network-online.target rsyslog.service
Requires=tink-agent-setup.service containerd.service
Wants=network-online.target rsyslog.service
# Only start when Tinkerbell parameters are on the kernel command line
ConditionKernelCommandLine=|tink_worker_image
ConditionKernelCommandLine=|docker_registry

[Service]
Type=simple
EnvironmentFile=/etc/tink-agent.env
ExecStart=/usr/local/bin/tink-agent-start

# Failure handling — restart with backoff
Restart=always
RestartSec=5
RestartMaxDelaySec=60
RestartSteps=5

# Logging to journal (forwarded to remote syslog by rsyslog)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tink-agent

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
